# templates/fsck/cvmfs-fsck.service.epp
<%-|
  String[1] $cache_base,
  Optional[String] $options = undef,
|-%>
# Puppet Installed

[Unit]
Description=cvmfs fsck
ConditionPathExists=<%= $cache_base %>/shared/txn

[Service]
Type=oneshot
Nice=19
IOSchedulingClass=2
IOSchedulingPriority=7
ExecStart=/usr/bin/cvmfs_fsck <%= $options %> <%= $cache_base %>/shared
